# crontab
*/5 * * * * /usr/bin/bash /root/gen-files.sh
