# crontab
*/1 * * * * /usr/bin/bash /root/gen-files.sh  >> /var/log/cron.log 2>&1
